# .config/autostart/electron.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=electron
Comment=electronstartup script
Exec=/usr/lib/electron/electron
StartupNotify=false
Terminal=false