# systemd/system/brightnessOn.service
[Unit]
Description=brightnessOn

[Service]
ExecStart=/usr/bin/brightnessOn

[Install]
WantedBy=multi-user.target
